# datasync.path
[Unit]
Wants=datasync.service

[Path]
PathChanged=/home/thomas/doc/
PathChanged=/home/thomas/img/
PathChanged=/home/thomas/_journal
PathChanged=/home/thomas/lmms/projects/
PathChanged=/home/thomas/do
PathChanged=/home/thomas/_poetry
PathChanged=/home/thomas/aud/
PathChanged=/home/thomas/.password-store
PathChanged=/home/thomas/vid
PathChanged=/home/thomas/dev
Unit=datasync.service

[Install]
WantedBy=datasync.service
